# .internal-ci/docker/support/node_hw/supervisor/conf.d/admin-http-gw.conf
; Copyright (c) 2018-2022 The MobileCoin Foundation
[program:mc-admin-http-gateway]
priority=100
; If we don't start in 20 seconds, then go to fatal.
startsecs=20
command=/usr/bin/mc-admin-http-gateway
    --listen-host 0.0.0.0
    --listen-port 8000
    --admin-uri insecure-mca://127.0.0.1:8001/

stdout_logfile=/dev/null
stderr_logfile=/dev/null
autorestart=true
