# .internal-ci/docker/support/node_hw/supervisor/conf.d/aesm-service.conf
; Copyright (c) 2018-2022 The MobileCoin Foundation
[program:aesm-service]
priority=50
; if we don't start in 20 seconds go to fatal
startsecs=20
environment=AESM_PATH="/opt/intel/sgx-aesm-service/aesm",LD_LIBRARY_PATH="/opt/intel/sgx-aesm-service/aesm"
command=/opt/intel/sgx-aesm-service/aesm/aesm_service
    --no-daemon
    --no-syslog

stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
