# .internal-ci/docker/support/node_hw/supervisor/conf.d/consensus-service.conf
; Copyright (c) 2018-2022 The MobileCoin Foundation
[program:consensus-service]
priority=100
; if we don't start in 60 seconds go to fatal.
startsecs=60
command=/usr/local/bin/wrapper-consensus-service.sh

stdout_syslog=true
stderr_logfile=true

stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
