# .internal-ci/docker/support/node_hw/supervisor/conf.d/ledger-distribution.conf
; Copyright (c) 2018-2022 The MobileCoin Foundation
[program:ledger-distribution]
priority=50
; if we don't start in 20 seconds go to fatal
startsecs=20
command=/usr/local/bin/wrapper-ledger-distribution.sh

stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
